# cloud/ansible/roles/frontend/templates/redis.container.j2
[Service]
TimeoutStartSec=1800

[Container]
Image=docker.io/redis/redis-stack:latest
ContainerName=redis
Network=ollama.network
