# cloud/ansible/roles/ollama/templates/ollama.container.j2
[Container]
Image=docker.io/ollama/ollama:{{ ollama_container_image_tag }}
ContainerName=ollama
Volume=ollama.volume:/root/.ollama
AddDevice=nvidia.com/gpu=all
SecurityLabelDisable=true
Network=ollama.network

[Install]
WantedBy=default.target