# cloud/ansible/roles/ollama/templates/slackbot.container.j2
[Unit]
Requires=ingest.service
After=ingest.service

[Service]
TimeoutStartSec=1800

[Container]
Image={{ ollama_client_container_image }}
ContainerName=slackbot
Volume=embedding.volume:/db
Environment=PERSIST_DIRECTORY=/db
Environment=SLACK_APP_TOKEN={{ollama_slack_app_token}}
Environment=SLACK_BOT_TOKEN={{ollama_slack_bot_token}}
Environment=OLLAMA_HOST=ollama
Network=ollama.network
